# deploy/jarvis-buytracker.service
[Unit]
Description=JARVIS Buy Tracker Bot
Documentation=https://github.com/lucid/jarvis
After=network-online.target
Wants=network-online.target
PartOf=jarvis.target
ReloadPropagatedFrom=jarvis.target

[Service]
Type=simple
User=jarvis
Group=jarvis
WorkingDirectory=/opt/jarvis
ExecStart=/usr/bin/python3 -m bots.buy_tracker.bot
Restart=always
RestartSec=15
TimeoutStopSec=30

# Environment
EnvironmentFile=/opt/jarvis/.env
Environment="BOT_COMPONENT=buytracker"
Environment="SKIP_SUPERVISOR=1"

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/jarvis/data /opt/jarvis/logs /opt/jarvis/.lifeos
PrivateTmp=true

# Resource limits
MemoryMax=512M
CPUQuota=30%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jarvis-buytracker

# Watchdog
WatchdogSec=600

[Install]
WantedBy=jarvis.target
WantedBy=multi-user.target
